# deploy/systemd/v5-api.service.tmpl
[Unit]
Description=V5 read-only dashboard / API
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__V5_HOME__
Environment=API_HOST=127.0.0.1
Environment=API_PORT=8765
EnvironmentFile=__V5_HOME__/.env
ExecStart=__V5_HOME__/.venv/bin/python __V5_HOME__/scripts/run_api.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target