# deploy/systemd/v5-health.service
[Unit]
Description=V5 health check (oneshot)
Wants=v5-health.timer

[Service]
Type=oneshot
Environment=V5_HOME=%h/v5
WorkingDirectory=${V5_HOME}
EnvironmentFile=${V5_HOME}/.env
ExecStart=${V5_HOME}/.venv/bin/python ${V5_HOME}/scripts/health_check.py \
    --services v5-trader v5-api \
    --api-url http://127.0.0.1:8765 \
    --format text
# Non-zero is logged but does not retry; the timer drives the cadence.
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
PrivateTmp=yes
