# deploy/systemd/v5-trader.service
[Unit]
Description=V5 paper-orchestrator (scanner + exec loop)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
# Reference unit for a default ~/v5 install under systemd --user.
# For production installs on arbitrary paths, prefer `sudo deploy/install.sh`.
WorkingDirectory=%h/v5
EnvironmentFile=%h/v5/.env
ExecStart=%h/v5/.venv/bin/python %h/v5/scripts/run_exec.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=%h/v5/data

[Install]
WantedBy=default.target