# deploy/systemd/v5-trader.service
[Unit]
Description=V5 paper-orchestrator (scanner + exec loop)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
# Reference unit for a default ~/v5 install under systemd --user.
# For production installs on arbitrary paths, prefer `sudo deploy/install.sh`
# which renders system-level units with absolute paths.
WorkingDirectory=%h/v5
EnvironmentFile=%h/v5/.env
ExecStart=%h/v5/.venv/bin/python %h/v5/scripts/run_exec.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=%h/v5/data

[Install]
WantedBy=default.target[Unit]
Description=V5 paper-orchestrator (scanner + exec loop)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
# Default install path is ~/v5. If your clone lives elsewhere (e.g. ~/V1)
# create a drop-in override:
#   systemctl --user edit v5-trader.service
# and override WorkingDirectory / EnvironmentFile / ExecStart together
# (Environment= variables do NOT expand in WorkingDirectory/EnvironmentFile).
WorkingDirectory=%h/v5
EnvironmentFile=%h/v5/.env
# Loopback-only by default; rely on the venv's Python and our launcher.
ExecStart=%h/v5/.venv/bin/python %h/v5/scripts/run_exec.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
# stdout/stderr go to journald; we still write structured JSON logs to
# data/logs/ for replay and grep.
StandardOutput=journal
StandardError=journal
# --- Hardening (minimal: most directives need root or unprivileged-userns) ---
# When you eventually move to a system unit, re-enable: NoNewPrivileges,
# PrivateTmp, ProtectSystem=full, ProtectHome, MemoryMax, LimitNOFILE,
# RestrictRealtime. Under --user systemd they all risk 218/CAPABILITIES
# or "resources" failures depending on the kernel + cgroup delegation.
NoNewPrivileges=yes

[Install]
WantedBy=default.target
